# templates/default/rc.local
sh /etc/init.d/networking restart # vagrant ssh hotfix
exit 0